# recipes-connectivity/wlan-enabler/wlan-enabler/wlan-enabler.service
[Unit]
Description=Bunch-Linux wlan enabler service

[Service]
Type=oneshot
ExecStart=/usr/bin/wlan-enabler.sh

[Install]
WantedBy=multi-user.target
