# recipes-app/systemdservices/buttonreaderservices.service
[Unit]
Description=Button reader service

[Service]
User=root
ExecStart=/usr/bin/zmqread

[Install]
WantedBy=multi-user.target
