# recipes-core/bunch-update/bunch-update/update.timer
[Unit]
Description=Update timer

[Timer]
OnBootSec=5min
Unit=update.service

[Install]
WantedBy=timers.target